# asterisk-integration.service
[Unit]
Description=Asterisk CRM Integration Service
Documentation=https://github.com/your-org/asterisk-integrations
After=network.target
# If Asterisk runs on the same machine, wait for it too:
# After=network.target asterisk.service

[Service]
Type=simple

# ── User ──────────────────────────────────────────────────────────────────────
# Run as root OR set to an existing user on this machine.
# To find a valid user, run:  cut -d: -f1 /etc/passwd | sort
# Common choices: root | asterisk | nobody
User=root

# ── Paths ─────────────────────────────────────────────────────────────────────
# Absolute path to your deployment directory.
WorkingDirectory=/var/www/asterisk-integrations

# PHP binary path — find yours with:  which php   OR   whereis php
# Common paths:
#   /usr/bin/php
#   /usr/local/bin/php
#   /usr/bin/php8.1
ExecStart=/usr/bin/env php /var/www/asterisk-integrations/service.php

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=always
RestartSec=5

# Give it 10 seconds to shut down gracefully before SIGKILL.
TimeoutStopSec=10

# ── Logging ───────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asterisk-integration

# ── Hardening ─────────────────────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
